# launchagents/com.hausverwaltung.coverage-audit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hausverwaltung.coverage-audit</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-c</string>
        <string>/Users/dan/.local/bin/hausverwaltung-coverage-audit.sh</string>
    </array>

    <!-- Sunday + Wednesday at 03:00 local time (2x/week, ~3.5 day spacing).
         Sunday 03:00 = end-of-week sweep; Wednesday 03:00 = mid-week catch.
         03:00 chosen so the user is asleep and runs finish long before
         work starts. -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Weekday</key>
            <integer>0</integer>
            <key>Hour</key>
            <integer>3</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>3</integer>
            <key>Hour</key>
            <integer>3</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/dan/Library/Logs/hausverwaltung-coverage-audit.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/dan/Library/Logs/hausverwaltung-coverage-audit.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
