# launchagents/com.hausverwaltung.deploy-validate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hausverwaltung.deploy-validate</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-c</string>
        <string>/Users/dan/.local/bin/hausverwaltung-deploy-validate.sh</string>
    </array>

    <!-- Daily 01:00 local. Runs before all other crons (cso 02:00, coverage 03:00)
         so the validate stack tests the *current* main, not main-after-tonight's-fixes. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>1</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/dan/Library/Logs/hausverwaltung-deploy-validate.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/dan/Library/Logs/hausverwaltung-deploy-validate.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
