# launchagents/com.hausverwaltung.openclaw-ui-watch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hausverwaltung.openclaw-ui-watch</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-c</string>
        <string>/Users/dan/.local/bin/hausverwaltung-openclaw-ui-watch.sh</string>
    </array>

    <!-- Weekly, Monday 02:00 local. Cheap-idle most weeks (git pull + a diff
         that touches no chat files → no-op, no claude run). On change weeks it
         runs a single read-only `claude -p` evaluation and may open ONE gated
         GitHub issue (label openclaw-ui-watch). Never writes code. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/dan/Library/Logs/hausverwaltung-openclaw-ui-watch.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/dan/Library/Logs/hausverwaltung-openclaw-ui-watch.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
